# launchd/com.cbmartinx.tidal-automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.cbmartinx.tidal-automation</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/chris/Documents/Scripts/tidal-automation/.venv/bin/python3</string>
      <string>-u</string>
      <string>/Users/chris/Documents/Scripts/tidal-automation/tidal_automation.py</string>
      <string>all</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/chris/Documents/Scripts/tidal-automation</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/chris/Documents/Scripts/tidal-automation/logs/tidal-automation.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chris/Documents/Scripts/tidal-automation/logs/tidal-automation.err.log</string>
  </dict>
</plist>
